# dist/upstart/horepgd.conf
description "Import EPG data from horizon.tv"
author      "Beralt Meppelink"

start on started tvheadend
stop on runlevel [016]

# Automatically restart process if crashed
respawn
expect daemon

script
  exec horepgd.py -d
end script
