# dist/systemd/horepgd.service
[Unit]
Description=Horizon EPG update service for TVHeadend
After=tvheadend.service
Requires=tvheadend.service
 
[Service]
Type=oneshot
Nice=19
ExecStart=/usr/bin/horepgd.py -u hts -g video -R -1

[Install]
WantedBy=multi-user.target
